# exim4/files/exim.init
#!/bin/sh /etc/rc.common
# Exim
 
START=88
STOP=12
 
start() {
    chown root:exim /usr/exim/configure
    /usr/exim/bin/exim -bd -q30m
}
 
stop() {
    killall exim
}

restart() {
    stop
    start
}

