# exim4/files/exim.init
#!/bin/sh /etc/rc.common
# Exim
 
START=91
STOP=09
PIDFILE=/var/run/exim.pid
 
start() {
    /usr/exim/bin/exim -bd -q30m
}
 
stop() {
    kill $(cat "$PIDFILE")
}

restart() {
    stop
    start
}

