# beanstalkd/files/beanstalkd.init
#!/bin/sh /etc/rc.common
# Beanstalk
 
START=89
STOP=11
 
start() {
    /usr/bin/beanstalkd &
}
 
stop() {
    killall beanstalkd
}

restart() {
    stop
    start
}

