# cfengine/files/cfengine.init
#!/bin/sh /etc/rc.common
# cfengine
 
START=90
STOP=10
 
start() {
    echo start
    mkdir -p /var/cfengine/bin
    cp /opt/cfengine/bin/cf-promises /var/cfengine/bin/cf-promises
    /opt/cfengine/bin/cf-key &> /dev/null
    /opt/cfengine/bootstrap.cf
    echo started.
}
 
