# cfengine/files/cfengine.init
#!/bin/sh /etc/rc.common
# Example script
# Copyright (C) 2007 OpenWrt.org
 
START=90
STOP=10
 
start() {
    echo start
    mkdir -p /var/cfengine/bin
    cp /opt/cfengine/bin/cf-promises /var/cfengine/bin/cf-promises
    /opt/cfengine/bin/cf-key
    /opt/enigmabox/cfengine/bootstrap.cf
    echo started.
}
 
stop() {
    echo stop
}
