# cfengine-promises/files/root/system_base/templates/crontab

# cjdns networking
* * * * * /usr/sbin/setup-cjdns-networking &> /dev/null

# subscriber stuff
0 */4 * * * /usr/sbin/subscriber-stuff &> /dev/null

