# cjdns-master/files/cjdns.init
#!/bin/sh /etc/rc.common
# cjdns
 
START=95
STOP=05
PIDFILE=/var/run/cjdns
EXTRA_COMMANDS="status"
 
start() {
    (/usr/sbin/cjdroute < /box/cjdroute.conf) &
}
 
stop() {
    killall cjdroute
}

status() {
    echo stopped
}

restart() {
    stop
    start
}

