# cjdns-master/files/cjdns.init
#!/bin/sh /etc/rc.common
# cjdns
 
START=93
STOP=07
EXTRA_COMMANDS="status"
 
start() {
    ( /usr/sbin/cjdroute < /box/cjdroute.conf &> /dev/null ) &
}
 
stop() {
    killall cjdroute
}

status() {
    if [[ $(ps | grep cjdroute | wc -l) -gt 0 ]]; then
        echo running
    else
        echo stopped
    fi
}

restart() {
    stop
    start
}

