# cfengine-promises/files/cfengine-promises.init
#!/bin/sh /etc/rc.common
# cfengine promises apply service

START=92

start() {
    /usr/sbin/cfengine-apply
    /usr/sbin/asterisk -rx 'sip notify gsreboot 100'
}

