# cfengine-policies/files/cfengine-policies.init
#!/bin/sh /etc/rc.common
# cfengine policies apply service
 
START=92
 
start() {
    /usr/sbin/cfengine-apply
}
 
