# cfengine-policies/files/rebuild-iptables.init
#!/bin/sh /etc/rc.common
# rebuild iptables
 
START=90
 
start() {
    /usr/sbin/rebuild-iptables
}
 
