# cfengine-policies/files/setup-cjdns-networking.init
#!/bin/sh /etc/rc.common
# setup cjdns networking
 
START=94
 
start() {
    /usr/sbin/setup-cjdns-networking
}
 
